# remna-agent.service
[Unit]
Description=Remna Node Agent (problem-first telemetry)
After=network-online.target docker.service
Wants=network-online.target

[Service]
Type=simple
User=remna-agent
Group=remna-agent
SupplementaryGroups=docker
WorkingDirectory=/opt/remna-agent
EnvironmentFile=/opt/remna-agent/.env
ExecStart=/opt/remna-agent/.venv/bin/python /opt/remna-agent/agent.py
Restart=always
RestartSec=3
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=true
ReadWritePaths=/opt/remna-agent /var/log/remnanode
ReadOnlyPaths=/var/log/remnanode

[Install]
WantedBy=multi-user.target
